# images/capi/ansible/roles/providers/files/etc/systemd/system/cloud-init.service.d/boot-order.conf
[Unit]
After=NetworkManager-wait-online.service